# ports/pasture/gtk+2/profile.d/ChangeLog
2005-04-08 12:54  topf

	* gtk+.csh, gtk+.sh: Added.

