# ports/testing/hal/rc.hald.new
#!/bin/sh
#
# haldaemon:   HAL daemon
#
# chkconfig: 345 98 02
# description:  This is a daemon for collecting and maintaing information \
#               about hardware from several sources. \
#               See http://www.freedesktop.org/software/hal/
#
# processname: hald
# 
#
# Sanity checks.
[ -x /usr/sbin/hald ] || exit 0

RETVAL=0

start() {
  echo "Starting HAL daemon:  hald"
  cat /var/run/hald.pid 2>/dev/null  && echo "HAL daemon is already running." && exit
  hald --retain-privileges
 #RETVAL=$?
 #[ $RETVAL -eq 0 ] && echo `/sbin/pidof hald` >/var/run/hald/pid
}

stop() {
  echo "Stopping HAL daemon..."
  kill -TERM `cat /var/run/hald.pid 2>/dev/null` 2>/dev/null || echo "HAL daemon was not running."
  rm -f /var/run/hald.pid
  #RETVAL=$?
  #[ $RETVAL -eq 0 ] && rm -f /var/run/hald/pid

}

# See how we were called.
case "$1" in
  start)
    start;;
  stop)
    stop;;
  restart)
    stop
    sleep 1
    start;;
  *)
    echo $"Usage: $0 {start|stop|restart}";;
esac
exit $RETVAL
